# vhclient/rootfs/etc/s6-overlay/s6-rc.d/vhclient-use-devices/up
/etc/s6-overlay/scripts/vhclient-use-devices