# vhclient/rootfs/etc/s6-overlay/s6-rc.d/vhclient/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting VirtualHere Client"
exec /usr/bin/vhclient -c /config/vhclient.ini -l /dev/stdout
